# update.cron.template
#run the update script every 4 hours to ensure mods are usually up to date.
00 04 * * * ${repo_dir}/update.sh -b ${config_branch}